# backend/cron/crontab
# Run scrape_meget every 6 hours (10 pages ~ 200 listings)
0 */6 * * * root cd /app && /usr/local/bin/python -m flask scrape_meget --pages 10 >> /var/log/cron.log 2>&1
# Run scrape_bon_ua every 6 hours (offset by 1 hour)
0 1-23/6 * * * root cd /app && /usr/local/bin/python -m flask scrape_bon_ua --pages 10 >> /var/log/cron.log 2>&1
