# backend/cron/crontab
# Run scrape every 6 hours (10 pages ~ 200 listings)
0 */6 * * * root cd /app && /usr/local/bin/python -m flask scrape_meget --pages 10 >> /var/log/cron.log 2>&1
# An empty line is required at the end of this file for a valid cron file.
