# backend/cron/crontab
# Run scrape every 6 hours
0 */6 * * * root cd /app && /usr/local/bin/python -m flask scrape >> /var/log/cron.log 2>&1
# An empty line is required at the end of this file for a valid cron file.
